# etc_ubuntu/config/autostart/gnome-terminal.desktop
[Desktop Entry]
Exec="gnome-terminal --geometry=1920*1080+0+0"
Type=Application
